# .config/autostart/pamac-tray.desktop
[Desktop Entry]
Exec=/usr/bin/pamac-tray
Icon=
Name=pamac-tray
Path=
Terminal=False
Type=Application
